# root/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

if [ $JACKETT_NO_UPDATES -eq 1 ]; then
  UPDATE="--NoUpdates"
else
  UPDATE=""
fi

exec \
	s6-setuidgid abc /app/Jackett/jackett_launcher.sh $UPDATE "${RUN_OPTS}"
